# docker/supervisor.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
logfile_maxbytes=50MB
logfile_backups=10

[program:laravel]
command=php artisan serve --host=0.0.0.0 --port=8000
directory=/var/www
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/laravel.log
stderr_logfile=/var/log/supervisor/laravel.err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=3600
